# scripts/com.pruviq.e2e-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pruviq.e2e-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jepo/pruviq/scripts/e2e-monitor.sh</string>
    </array>

    <!-- Every 6 hours: 00:30, 06:30, 12:30, 18:30 UTC -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>3</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/pruviq-e2e/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/pruviq-e2e/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/opt/node@22/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/jepo</string>
    </dict>

    <!-- Don't restart on failure -->
    <key>KeepAlive</key>
    <false/>

    <!-- 10 min timeout -->
    <key>TimeOut</key>
    <integer>600</integer>
</dict>
</plist>
